# stow/kb-kill/.config/systemd/user/kb-kill.service
[Unit]
Description=kb-kill - disable/enable a target keyboard on a global hotkey
Documentation=https://github.com/seisinger/dotfiles
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/kb-kill run
Restart=on-failure
RestartSec=2
# Hardening: this only needs to read/write /dev/input/event* (input group).
NoNewPrivileges=true

[Install]
WantedBy=graphical-session.target
