# scripts/kb-kill/kb-kill-tray.service
[Unit]
Description=kb-kill-tray - tray icon for kb-kill
Documentation=https://github.com/seisinger/dotfiles
# kb-kill itself is a system service; the tray reconnects on its own, so no
# ordering dependency on it is needed here.
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/kb-kill-tray
Restart=on-failure
RestartSec=3

[Install]
WantedBy=graphical-session.target
